# boxes/demos/gptvswiki/daemon/chromium.timer
[Unit]
Description=Chromium Daemon (timer)

[Timer]
OnBootSec=5sec
OnUnitActiveSec=10sec

[Install]
WantedBy=timers.target