# boxes/demos/gptvswiki/daemon/gptvswiki.timer
[Unit]
Description=GPTvsWiki Daemon (timer)

[Timer]
OnBootSec=1sec
OnUnitActiveSec=10sec

[Install]
WantedBy=timers.target